# conf/systemd.service
[Unit]
Description=Docker container: __IMAGE__ - YunoHost app __APP__
Requires=docker.service
After=docker.service

[Service]
Type=simple
User=__APP__
Group=__APP__

# Remove any leftover container from a previous run before starting
ExecStartPre=-/usr/bin/docker rm -f __APP__

# Run the container in the foreground so systemd tracks it.
# Port mapping is automatic: host port __PORT__ -> container port __CONTAINER_PORT__
# Note: --rm is incompatible with --restart, so we manage cleanup in ExecStop instead
ExecStart=/usr/bin/docker run \
    --restart=__DOCKER_RESTART_POLICY__ \
    --name=__APP__ \
    --log-driver=journald \
    --log-opt tag=__APP__ \
    -e TZ=UTC \
    --volume /etc/timezone:/etc/timezone:ro \
    -p 127.0.0.1:__PORT__:__CONTAINER_PORT__ \
    __DOCKER_VOLUME_OPTS__ \
    __DOCKER_OPTIONS__ \
    __IMAGE__

ExecStop=-/usr/bin/docker stop __APP__
ExecStopPost=-/usr/bin/docker rm -f __APP__

Restart=no
RestartSec=10

### Security hardening
### Note: PrivateDevices, RestrictNamespaces, NoNewPrivileges, SystemCallFilter are
### intentionally omitted — they are incompatible with running Docker containers.
### Running as a dedicated system user in the docker group is the primary mitigation.
PrivateTmp=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
