# conf/systemd.service
[Unit]
Description=Docker container: __IMAGE_REF__ - YunoHost app __APP__
Requires=docker.service
After=docker.service

[Service]
Type=simple
User=__APP__
Group=__APP__

# Remove any leftover container from a previous run before starting
ExecStartPre=-/usr/bin/docker rm -f __APP__

# Run the container in the foreground so systemd tracks it.
# Port mapping is automatic: host port __PORT__ -> container port __CONTAINER_PORT__
# Note: --rm is incompatible with --restart, so we manage cleanup in ExecStop instead
ExecStart=/usr/bin/docker run \
    --restart=__DOCKER_RESTART_POLICY__ \
    --name=__APP__ \
    --log-driver=journald \
    --log-opt tag=__APP__ \
    --security-opt=no-new-privileges \
    -e TZ=UTC \
    --volume /etc/timezone:/etc/timezone:ro \
    -p 127.0.0.1:__PORT__:__CONTAINER_PORT__ \
    __DOCKER_NETWORK_OPT__ \
    __DOCKER_VOLUME_OPTS__ \
    __DOCKER_OPTIONS__ \
    __IMAGE_REF__

ExecStop=-/usr/bin/docker stop __APP__
ExecStopPost=-/usr/bin/docker rm -f __APP__

Restart=no
RestartSec=10

### Security hardening
### This unit runs the `docker` *client*, which only connects to the Docker
### daemon over /var/run/docker.sock; dockerd creates and runs the container
### itself, so sandboxing this process does not constrain the container. We can
### therefore apply the full example_ynh hardening baseline to the client.
### Baseline: https://github.com/YunoHost/example_ynh/blob/main/conf/systemd.service
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
DevicePolicy=closed
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProtectSystem=full
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=~@clock @debug @module @mount @obsolete @reboot @setuid @swap @cpu-emulation @privileged
CapabilityBoundingSet=~CAP_RAWIO CAP_MKNOD
CapabilityBoundingSet=~CAP_AUDIT_CONTROL CAP_AUDIT_READ CAP_AUDIT_WRITE
CapabilityBoundingSet=~CAP_SYS_BOOT CAP_SYS_TIME CAP_SYS_MODULE CAP_SYS_PACCT
CapabilityBoundingSet=~CAP_LEASE CAP_LINUX_IMMUTABLE CAP_IPC_LOCK
CapabilityBoundingSet=~CAP_BLOCK_SUSPEND CAP_WAKE_ALARM
CapabilityBoundingSet=~CAP_SYS_TTY_CONFIG
CapabilityBoundingSet=~CAP_MAC_ADMIN CAP_MAC_OVERRIDE
CapabilityBoundingSet=~CAP_NET_ADMIN CAP_NET_BROADCAST CAP_NET_RAW
CapabilityBoundingSet=~CAP_SYS_ADMIN CAP_SYS_PTRACE CAP_SYSLOG

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
